# supervisord.conf
[program:node]
command = socketstream start
directory = /home/dotcloud/current
